# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

cd /app/radarr/bin || exit

exec \
	s6-setuidgid abc /app/radarr/bin/Radarr \
	-nobrowser -data=/config