# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

cd /opt/radarr || exit

exec \
	s6-setuidgid abc mono --debug Radarr.exe \
	-nobrowser -data=/config
