# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

cd /opt/radarr || exit

exec \
	s6-setuidgid abc mono Radarr.exe \
	-nobrowser -data=/config
